# util/MsSql/crontab
50 23 * * * bitwarden /backup-db.sh >>/var/log/cron.log 2>&1

# An empty line is required at the end of this file for a valid cron file.
