# docker-unified/supervisord/scim.ini
[program:scim]
autostart=true
autorestart=true
command=/usr/bin/dotnet "Scim.dll"
directory=/app/Scim
environment=ASPNETCORE_URLS="http://+:5002"
redirect_stderr=true
startsecs=15
stdout_logfile=/var/log/bitwarden/scim.log
