# src/Jobs/crontab
* * * * * root dotnet /jobs/Jobs.dll -d /jobs -j hello >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.